# ansible/roles/metrics-exporter/templates/golang-app.service.j2
[Unit]
Description=Go Lang Application
After=network.target

[Service]
ExecStart=/opt/go-app/main
Restart=always
User=root
EnvironmentFile=/etc/app.vars

[Install]
WantedBy=multi-user.target