# ansible/roles/simple-app/templates/golang-app.service.j2
[Unit]
Description=Go Lang Application
After=network.target

[Service]
ExecStart=/opt/go-app/main
Restart=always
EnvironmentFile=/etc/app.vars
User=root

[Install]
WantedBy=multi-user.target